# supervisor/apache.conf
[program:apache2]
command=chown -R www-data:www-data /srv/reviews.local/data && /bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
autorestart=true