# supervisor/apache.conf
[program:apache2]
command=service apache2 start
autorestart=false